# res/lapas/guest/etc/systemd/system/lapas-api-daemon.service
[Unit]
After=network.target

[Service]
Type=simple
ExecStart=/lapas/lapas-api-client daemon

[Install]
WantedBy=multi-user.target