# res/server/etc/systemd/system/lapas-api-server.service
[Unit]
After=network.target

[Service]
Type=simple
ExecStart=@@LAPAS_SCRIPTS_DIR@@/apiServer.sh

[Install]
WantedBy=multi-user.target